# deploy/evidence-bench-compose.service.example
[Unit]
Description=Evidence Bench shared laboratory service
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=INSTALL_DIRECTORY
ExecStart=/usr/bin/docker compose -f compose.yaml -f compose.local.yaml up -d --no-build --remove-orphans
ExecStartPost=/usr/bin/docker compose -f compose.yaml -f compose.local.yaml ps
ExecStop=/usr/bin/docker compose -f compose.yaml -f compose.local.yaml stop
TimeoutStartSec=900
TimeoutStopSec=180

[Install]
WantedBy=default.target
